# ebaseline/airootfs/etc/systemd/system/run-cryptsetup\x2dkeys.d.automount
[Unit]
Description=Keys directory bind automount

[Automount]
Where=/run/cryptsetup-keys.d
TimeoutIdleSec=3min

[Install]
WantedBy=multi-user.target
