# ebaseline/airootfs/etc/systemd/system/run-archiso-donglepersistent.automount
[Unit]
Description=Dongle persistent partition automount

[Automount]
Where=/run/archiso/donglepersistent
TimeoutIdleSec=3min

[Install]
WantedBy=multi-user.target
