# ebaseline/airootfs/etc/systemd/system/run-archiso-persistent.automount
[Unit]
Description=Persistent partition automount

[Automount]
Where=/run/archiso/persistent
TimeoutIdleSec=3min

[Install]
WantedBy=multi-user.target
